# modules/launchagents/com.alex.macos-session-health.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.alex.macos-session-health</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Users/alex/dotfiles/modules/bin/macos-session-health</string>
    <string>--quiet</string>
    <string>--top</string>
    <string>40</string>
    <string>--fd-warn</string>
    <string>512</string>
    <string>--codex-app</string>
    <string>/Applications/Codex.app</string>
    <string>--codex-config</string>
    <string>/Users/alex/.codex/config.toml</string>
    <string>--codex-chrome-native-hosts</string>
    <string>/Users/alex/.codex/chrome-native-hosts-v2.json</string>
    <string>--codex-node-repl-warn</string>
    <string>20</string>
    <string>--syspolicyd-fd-warn</string>
    <string>2048</string>
    <string>--syspolicyd-path-warn</string>
    <string>512</string>
    <string>--syspolicyd-lsof-interval-minutes</string>
    <string>5</string>
    <string>--syspolicyd-cpu-warn</string>
    <string>80</string>
    <string>--syspolicyd-rss-warn-mb</string>
    <string>1536</string>
    <string>--syspolicyd-rss-error-mb</string>
    <string>2048</string>
    <string>--syspolicyd-rss-growth-warn-mb-per-minute</string>
    <string>128</string>
    <string>--syspolicyd-log-error-count</string>
    <string>100</string>
    <string>--trustd-cpu-warn</string>
    <string>80</string>
    <string>--trustd-rss-warn-mb</string>
    <string>1024</string>
    <string>--trustd-rss-error-mb</string>
    <string>2048</string>
    <string>--trustd-rss-growth-warn-mb-per-minute</string>
    <string>128</string>
    <string>--code-sign-clone-warn-mb</string>
    <string>1024</string>
    <string>--code-sign-clone-growth-warn-mb-per-minute</string>
    <string>512</string>
    <string>--code-sign-clone-timeout</string>
    <string>3</string>
    <string>--codex-trusted-root-check-interval-minutes</string>
    <string>60</string>
    <string>--codex-trusted-root-check-timeout</string>
    <string>3</string>
    <string>--lsof-timeout</string>
    <string>5</string>
    <string>--fd-top-interval-minutes</string>
    <string>5</string>
    <string>--app-assess-timeout</string>
    <string>5</string>
    <string>--spawn-check-timeout</string>
    <string>5</string>
    <string>--spctl-interval-minutes</string>
    <string>0</string>
    <string>--codesign-interval-minutes</string>
    <string>0</string>
    <string>--log-timeout</string>
    <string>5</string>
    <string>--log-lines</string>
    <string>120</string>
    <string>--passive-log-interval-minutes</string>
    <string>5</string>
    <string>--passive-log-last-minutes</string>
    <string>6</string>
    <string>--passive-log-lines</string>
    <string>80</string>
    <string>--codex-passive-log-interval-minutes</string>
    <string>30</string>
    <string>--log-excerpt-cooldown-minutes</string>
    <string>30</string>
    <string>--retention-days</string>
    <string>14</string>
    <string>--brrr-helper</string>
    <string>/Users/alex/.agents/skills/skills-stable__brrr-now/scripts/brrr-send.sh</string>
    <string>--brrr-notify-cooldown-minutes</string>
    <string>10</string>
    <string>--brrr-interruption-level</string>
    <string>passive</string>
    <string>--app</string>
    <string>/Applications/Visual Studio Code.app</string>
    <string>--app</string>
    <string>/Applications/Codex.app</string>
    <string>watch</string>
    <string>--interval</string>
    <string>60</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/Users/alex/Library/Logs/macos-session-health/agent.out.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/alex/Library/Logs/macos-session-health/agent.err.log</string>

  <key>WorkingDirectory</key>
  <string>/Users/alex/dotfiles</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/Users/alex/.local/share/mise/shims:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    <key>PYTHONUNBUFFERED</key>
    <string>1</string>
  </dict>
</dict>
</plist>
